# custom/manjaro-sway-settings/source/etc/profile.d/termite.sh
#!/bin/sh

export TERM=xterm-termite
export TERMINAL_COMMAND='/usr/bin/termite --exec'